# crontab
0 1 * * * /usr/local/bin/python /app/src/main.py >> /app/cron.log 2>&1
